# tools/gpkg-publisher/systemd/gpkg-publisher.service
[Unit]
Description=Build and publish GeminiOS GPKG repositories
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/gpkg-publisher.env
WorkingDirectory=/srv/geminios
ExecStart=/usr/bin/env python3 /srv/geminios/tools/gpkg-publisher/publish.py
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
