# tools/gpkg-publisher/systemd/gpkg-publisher.timer
[Unit]
Description=Run the GeminiOS GPKG bulk publisher on a schedule

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
RandomizedDelaySec=15m
Persistent=true
Unit=gpkg-publisher.service

[Install]
WantedBy=timers.target
